# magen_srv/etc/crontab
# m h  dom mon dow   command
0 5  * * *  /root/orefmap/server/cron_job.sh