# core/etc/crontab
# m h  dom mon dow   command
0 5  * * *  /root/magen/server/cron_job.sh